# pi_camera/etc/systemd/system/imager.timer
[Unit]
Description=Runs imager.service every 2 minutes

[Timer]
OnCalendar=*:*:0/30
AccuracySec=0.1

[Install]
WantedBy=basic.target
