# pi_camera/etc/systemd/system/upload.service
[Unit]
Description=Takes an image at every 5 minute step

[Service]
Type=oneshot
User=pi
ExecStart=/bin/bash /usr/local/bin/uploader.sh
