# pi_camera/DEBIAN/postinst
#!/bin/bash
chown root:root /etc/systemd/system/imager.*
systemctl start imager.timer
systemctl enable imager.timer

